# systemd/backend.service
[Unit]
Description=Payments Backend Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=appuser
WorkingDirectory=/srv/app
EnvironmentFile=/srv/app/.env
Environment=NODE_ENV=production
ExecStart=/srv/app/scripts/start-backend-safe.sh
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/app/logs

[Install]
WantedBy=multi-user.target

# ==============================================================================
# Deployment Schritte
# ==============================================================================
#
# 1. Erstelle Benutzer (falls nicht vorhanden):
#    sudo useradd -r -s /bin/false appuser
#
# 2. Kopiere Artefakte:
#    sudo mkdir -p /srv/app
#    sudo cp -r dist/ package*.json migrations/ scripts/ /srv/app/
#
# 3. Erstelle .env:
#    sudo nano /srv/app/.env
#
# 4. Setze Berechtigungen:
#    sudo chown -R appuser:appuser /srv/app
#    sudo chmod 600 /srv/app/.env
#    sudo chmod 700 /srv/app/scripts/start-backend-safe.sh
#
# 5. Installiere Service:
#    sudo cp systemd/backend.service /etc/systemd/system/backend.service
#    sudo systemctl daemon-reload
#    sudo systemctl enable --now backend.service
#
# 6. Logs:
#    sudo journalctl -u backend.service -f
#    sudo journalctl -u backend.service --since "1 hour ago"
#
# ==============================================================================
# Troubleshooting
# ==============================================================================
#
# Status prüfen:
#   sudo systemctl status backend.service
#
# Service neustarten:
#   sudo systemctl restart backend.service
#
# Logs der letzten 100 Zeilen:
#   sudo journalctl -u backend.service -n 100
#
# Live-Logs:
#   sudo journalctl -u backend.service -f
